# plugins/gluster-proc.crontab
*/1 * * * * root /usr/lib64/nagios/plugins/gluster/check_proc_status.py
